# result_vid_demo/week3/meta-smartmirror/recipes-apps/camera-manager/files/camera-manager.service
[Unit]
Description=Camera Manager (SmartMirror)
After=network.target modules-load.target


[Service]
Type=simple
ExecStart=/usr/bin/camera-manager
Restart=on-failure
Nice=-5
# Use realtime scheduler if desired (requires capabilities in image)
# ExecStartPre=/usr/bin/setcap cap_sys_nice+ep /usr/bin/camera-manager


[Install]
WantedBy=multi-user.target