# result_vid_demo/week3/meta-smartmirror/recipes-core/systemd/files/serial-getty@ttyAMA0.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root - 115200,38400,9600 ttyAMA0
Type=idle
Restart=always
